# config/systemd/user/arch-maid.timer
[Unit]
Description=arch-maid cleans every week

[Timer]
OnBootSec=5min
OnUnitActiveSec=1w
Persistent=true

[Install]
WantedBy=timers.target
